# deploy/systemd-oas-atlas.service
[Unit]
Description=OAS Atlas API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/oas-atlas
Environment=OAS_ATLAS_DB=/var/lib/oas-atlas/oas_atlas.db
ExecStart=/opt/oas-atlas/.venv/bin/oas-atlas --db /var/lib/oas-atlas/oas_atlas.db serve --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
User=oas-atlas
Group=oas-atlas

[Install]
WantedBy=multi-user.target
